# systemd/segecha.service
[Unit]
Description=Segecha Logistics Web Application
After=network.target

[Service]
User=segecha
Group=segecha
WorkingDirectory=/path/to/your/app
Environment="PATH=/path/to/your/app/venv/bin"
Environment="FLASK_APP=app.py"
Environment="FLASK_ENV=production"
ExecStart=/path/to/your/app/venv/bin/gunicorn --workers 3 --bind 127.0.0.1:8000 app:app
Restart=always

[Install]
WantedBy=multi-user.target 